# TEMP/data/etc/init.d/fan_commander
#!/bin/sh /etc/rc.common
START=99
start() {
    echo "1" > /sys/class/gpio/fan-hw/value 2>/dev/null
    /usr/bin/fan_commander.sh &
}
stop() {
    kill $(pgrep -f "fan_commander.sh") 2>/dev/null
}
